# deploy/daemon/breadbox.service.tmpl
[Unit]
Description=Breadbox (Docker Compose)
Documentation=https://github.com/canalesb93/breadbox
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__INSTALL_DIR__
# Preserve the caddy profile when DOMAIN is configured.
ExecStart=/bin/sh -c '__COMPOSE_CMD__ up -d'
ExecStop=/bin/sh -c '__COMPOSE_CMD__ down'
# Restart-on-failure doesn't apply to oneshot; the Docker daemon restarts
# containers via `restart: unless-stopped` in docker-compose.prod.yml.
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
